# src/pi-gen-stage2-luaml-3dposclient.d/stage2-luaml-3dposclient/02-3dpos-client/files/3dprinteros.service
[Unit]
Description=3DPrinterOS Client
After=network.target

[Service]
WorkingDirectory=/opt/3dprinteros-client/
ExecStart=/usr/bin/python2 launcher.py
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=always

[Install]
WantedBy=multi-user.target